# simple-file-server.conf
start on runlevel [2345]
stop on runlevel [!2345]

exec /etc/simple-file-server/simple_file_server.py
